# cron.d/process-graylog-notification
MAILTO=jost@mpdl.mpg.de

*/5 * * * * root /root/graylog-notification/process-graylog-notification.sh -d /var/www/local/json-queue -g /root/graylog-notification/it-mpdl.github.io -e evil_IPs.txt > /dev/null

