# cron.d/process-graylog-notification
MAILTO=jost@mpdl.mpg.de

*/5 * * * * root /root/graylog-notification/process-graylog-notification.sh -d /var/www/local/json-queue -g /root/graylog-notification/testgraylog -e evilIPs.lst > /dev/null

